# xorg/init.d/xconfig
#!/sbin/openrc-run

command=/usr/local/bin/xconfig
pidfile=/var/run/xconfig.pid
name="xconfig"

depend()
{
	provide xconfig
	after LOGIN cleanvar moused syscons dbus
	before gdm kdm lightdm pcdm sddm slim xdm
}

