# extra/ghostbsd/rc.local
#!/bin/sh
/usr/local/bin/python /usr/local/etc/wncc/wncc.py
